# arch/systemd-user/jarvis-mcp-register.service
[Unit]
Description=Register jarvis-linux-mcp as MCP server in IronClaw
After=default.target

[Service]
Type=oneshot
RemainAfterExit=true
Environment="DATABASE_BACKEND=libsql"
Environment="DATABASE_URL=unused://libsql"
Environment="LIBSQL_PATH=%h/.ironclaw/jarvis.db"
ExecStartPre=/usr/bin/mkdir -p %h/.ironclaw
ExecStartPre=/usr/bin/sh -c '[ -f %h/.ironclaw/.env ] || /usr/bin/cp /opt/jarvis-os/arch/templates/ironclaw.env %h/.ironclaw/.env'
ExecStart=/usr/bin/ironclaw mcp add jarvis-linux \
  --transport stdio \
  --command /usr/local/bin/jarvis-linux-mcp \
  --arg mcp-server

[Install]
WantedBy=default.target
