# arch/systemd-user/jarvis-voice-daemon.service
[Unit]
Description=jarvis-os voice daemon (ElevenLabs Conversational AI client)
After=pipewire.service
Wants=pipewire.service

[Service]
Type=simple
EnvironmentFile=-%h/.ironclaw/.env
ExecStart=/usr/local/bin/jarvis-voice-daemon
Restart=on-failure
RestartSec=3
# Logs JSON estructurados van a journald: journalctl --user -u jarvis-voice-daemon -f

[Install]
WantedBy=default.target
