# arch/systemd-user/jarvis-voice-daemon.service
[Unit]
Description=jarvis-os voice daemon (wake word + VAD + STT)
After=pipewire.service
Wants=pipewire.service

[Service]
Type=simple
WorkingDirectory=/opt/jarvis-os/voice_daemon
ExecStart=/usr/bin/uv run voice-daemon
Restart=on-failure
RestartSec=3
# Logs por stdout van a journald: journalctl --user -u jarvis-voice-daemon -f

[Install]
WantedBy=default.target
